# device-dovecot-el9/src/main/resources/usr/lib/systemd/system/dovecot-postdiscovery.target
# systemd service for managing all dovecot postconfiguration.

[Unit]
Description=Dovecot Configuration Postdiscovery
After=dovecot-autodiscovery.target

[Install]
WantedBy=multi-user.target

