# device-dovecot-el9/src/main/resources/usr/lib/systemd/system/dovecot-autodiscovery.target
# systemd service for managing all dovecot autoconfiguration.

[Unit]
Description=Dovecot Configuration Autodiscovery

[Install]
WantedBy=multi-user.target

